# dummy-data-provider/crontab
0 * * * * docker exec abacws-api node generate > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
